# roles/lvm/files/mantl-storage-setup.conf
# mountmantl - Mount cogscale filesystems on boot
#
# This helper mounts filesystems in the correct order as the devices
# and mountpoints become available.

description     "Mount cogscale filesystems on boot"

start on startup
stop on starting rcS

expect daemon
task

script
    exec /usr/bin/mantl-storage-setup
end script
